# backend/sausage-store-backend.service
[Unit]
Description=Sausage store
After=syslog.target network.target

[Service]
SuccessExitStatus=143

User=jarservice
Type=simple

<<<<<<< HEAD
<<<<<<< HEAD
=======
>>>>>>> 344506ec48f9ab7777af8a9c9ea885a9a13be3a8
StandardError=append:/opt/log/sausage-store.log
StandardOutput=append:/opt/log/sausage-store.log
Environment="REPORT_PATH=/opt/log/reports/"
Environment="LOG_PATH=/opt/log/"
Restart=always
Environment="LOGGING_LEVEL_ROOT=INFO"

#Start backend jar file
ExecStart=/usr/bin/java \
-jar '/opt/sausage-store.jar' \
--spring.datasource.url=jdbc:postgresql://c-c9qslvj961omn1nevnjl.rw.mdb.yandexcloud.net:6432/testdb \
--spring.datasource.username=sa \
--spring.datasource.password=password \
--spring.flyway.baselineOnMigrate=true

[Install]
WantedBy=multi-user.target


#ExecStart=java -jar /opt/sausage-store.jar

